# systemd/jellysink.service
[Unit]
Description=Jellysink media library scan service
After=network.target

[Service]
Type=oneshot
User=%u
ExecStart=/usr/local/bin/jellysinkd

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/jellysink %h/.config/jellysink

[Install]
WantedBy=default.target
