# systemd/jellysink.service
[Unit]
Description=Jellysink media library scan service
After=graphical-session.target
Documentation=https://github.com/Nomadcxx/jellysink

[Service]
Type=oneshot
User=%u
ExecStart=/usr/local/bin/jellysinkd

# Environment variables for GUI launching
Environment=DISPLAY=%E{DISPLAY}
Environment=WAYLAND_DISPLAY=%E{WAYLAND_DISPLAY}
Environment=XDG_RUNTIME_DIR=%t

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/jellysink %h/.config/jellysink

[Install]
WantedBy=graphical-session.target
